# grok-voice/grok-voice.service
[Unit]
Description=Grok Voice Stream
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/grok-voice
ExecStart=/usr/bin/python3 /home/ubuntu/grok-voice/grok_stream.py
Restart=always
RestartSec=5
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin

[Install]
WantedBy=multi-user.target
